# cron/monitor.crontab
*/5 * * * * PYTHONWARNINGS=ignore /usr/local/bin/python /app/monitor.py >> /proc/1/fd/1 2>&1
