# deploy/rakuten-bot.service
[Unit]
Description=Rakuten Monitor Discord Bot
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=yang_server
Group=yang_server
WorkingDirectory=/home/yang_server/rakuten
EnvironmentFile=-/home/yang_server/.rakuten_env
ExecStart=/usr/bin/python3 -m discord_bot
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rakuten-bot

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target